# ua-confession-mirror.cron
# Change the path to run.sh to your absolute path to the run.sh file
# This lines should be added when doing "crontab -e" in a terminal
# I didn't search how I could get this to work in a venv, so I had to install the requirements globally
#  Update: if you want you can activate the venv in run.sh, so you don't have to install the requirements globally

# run script every 60 minutes
*/60 * * * * /home/arno/Documents/ua-confessions-mirror/run.sh

# I haven't found an easy way to debug this, currently I set it to run every minute when debugging