# ua-confession-mirror.cron
# src: https://stackoverflow.com/questions/21648410/write-python-script-that-is-executed-every-5-minutes
# Change arno to your username, the path to main.py to your absolute path to the main.py file
# This two lines should be added when doing "crontab -e" in a terminal
# I didn't search how I could get this to work in a venv, so I had to install the requirements globally

# run script every 45 minutes
*/45 * * * *  python3 /home/arno/Documents/ua-confessions-mirror/main.py >>/home/arno/ua-confession-mirror-info.log 2>>/home/arno/ua-confession-mirror-error.log
# run script after system (re)boot
#@reboot python3 /home/arno/Documents/ua-confessions-mirror/main.py >>/home/arno/ua-confession-mirror-info.log 2>>/home/arno/ua-confession-mirror-error.log